# root/etc/services.d/apache/run
#!/usr/bin/with-contenv sh

# Start apache httpd.
exec /usr/sbin/httpd -DFOREGROUND
#sleep 3600
